# test/common_files/correct_test_result_json/filelist_info.json
{
  "code": "2001",
  "data": {
    "Judy": {
      "dir": [
        "/usr/share/licenses/Judy"
      ],
      "file": [
        "/usr/share/licenses/Judy/COPYING",
        "/usr/lib64/libJudy.so.1",
        "/usr/lib64/libJudy.so.1.0.3"
      ],
      "ghost": []
    },
    "dnf": {
      "dir": [
        "/var/cache/dnf",
        "/etc/bash_completion.d/dnf",
        "/etc/dnf",
        "/usr/share/doc/dnf",
        "/var/lib/dnf",
        "/var/lib/dnf/history",
        "/var/lib/dnf/yumdb",
        "/etc/dnf/aliases.d",
        "/etc/dnf/modules.d",
        "/etc/dnf/modules.defaults.d",
        "/etc/dnf/plugins",
        "/etc/dnf/protected.d",
        "/etc/dnf/vars",
        "/usr/share/licenses/dnf"
      ],
      "file": [
        "/usr/lib/systemd/system/dnf-automatic-download.service",
        "/usr/lib/systemd/system/dnf-automatic-download.timer",
        "/usr/lib/systemd/system/dnf-automatic-install.service",
        "/usr/lib/systemd/system/dnf-automatic-install.timer",
        "/usr/lib/systemd/system/dnf-automatic-notifyonly.service",
        "/usr/lib/systemd/system/dnf-automatic-notifyonly.timer",
        "/usr/lib/systemd/system/dnf-automatic.service",
        "/usr/lib/systemd/system/dnf-automatic.timer",
        "/usr/lib/systemd/system/dnf-makecache.service",
        "/usr/lib/systemd/system/dnf-makecache.timer",
        "/etc/dnf/protected.d/dnf.conf",
        "/usr/bin/dnf",
        "/usr/bin/dnf-automatic",
        "/usr/lib/tmpfiles.d/dnf.conf",
        "/usr/share/licenses/dnf/COPYING",
        "/usr/share/licenses/dnf/PACKAGE-LICENSING",
        "/etc/libreport/events.d/collect_dnf.conf",
        "/etc/logrotate.d/dnf",
        "/usr/share/doc/dnf/AUTHORS",
        "/usr/share/doc/dnf/README.rst",
        "/etc/dnf/automatic.conf",
        "/etc/dnf/dnf.conf"
      ],
      "ghost": [
        "/var/log/dnf.librepo.log",
        "/var/log/dnf.log",
        "/var/log/dnf.plugin.log",
        "/var/log/dnf.rpm.log",
        "/var/log/hawkey.log",
        "/var/lib/dnf/groups.json"
      ]
    }
  },
  "msg": "Successful Operation!"
}